# packaging.in/iot-launch.service.in
[Unit]
Description=IoT Application Launcher Daemon

[Service]
Type=simple
ExecStart=/usr/bin/iot-launch-daemon -A @libexecdir@/iot-launcher/io-launch-agent -f -vvv -S launcher,appfw
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
