# packaging.in/iot-launch.service.in
[Unit]
Description=IoT Application Launcher Daemon

[Service]
Type=simple
ExecStart=/usr/bin/iot-launch-daemon -a @libexecdir@/iot-launcher/iot-launch-agent -f -vvv
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
