# packaging.in/iotpm-register-apps.service.in
[Unit]
Description=registers the pre-installed applications with Security Manager
Requires=security-manager.socket
After=security-manager.socket

[Service]
type=oneshot
RemainAfterExit=no
TimeoutSec=0
ExecStart=@bindir@/register-preinstalled-apps

[Install]
WantedBy=multi-user.target
